# ipar.service
[Unit]
Description=Xnote IP address recorder

[Service]
Type=oneshot
ExecStart=/home/user/bin/ipar.sh
